# rc.local
#!/bin/bash



if [ -f /aafirstboot ]; then /aafirstboot start ; fi
/sbin/ifconfig wlan0 192.168.1.1
/bin/systemctl restart hostapd
/bin/systemctl restart dnsmasq
exit 0